# print-server/setup/99-usb-printer.rules
# udev-regel for automatisk printer-detect
# usb_interface-subsystem er pålitelig for bInterfaceClass-matching (ikke usb/bDeviceClass)
ACTION=="add", SUBSYSTEM=="usb_interface", ATTR{bInterfaceClass}=="07", \
    RUN+="/bin/systemd-run --no-block /usr/local/bin/add-printer.sh"
