# print-server/setup/first-boot.service
[Unit]
Description=Rock 3C first-boot provisioning
After=armbian-firstrun.service network-online.target
Wants=network-online.target
ConditionPathExists=!/.first-boot-done

[Service]
Type=oneshot
ExecStart=/usr/local/bin/first-boot.sh
RemainAfterExit=yes
StandardOutput=append:/var/log/first-boot.log
StandardError=append:/var/log/first-boot.log

[Install]
WantedBy=multi-user.target
